# systemd/cockpitlogin.service
[Unit]
Description=Cockpitlogin Service
Requires=cockpitlogin.socket

[Service]
ExecStart=/opt/cockpitlogin

