# systemd/cockpitlogin.socket
[Unit]
Description=Cockpitlogin Service Socket

[Socket]
ListenStream=/run/cockpitlogin/socket

[Install]
WantedBy=sockets.target
