# cockpitlogin/systemd/cockpitlogin.service
[Unit]
Description=Cockpitlogin Service
Requires=cockpitlogin.socket

[Service]
ExecStart=/usr/bin/cockpitlogin
RuntimeDirectory=cockpitlogin
